# slowcontrol/xsphere-slowcontrol.service
# systemd service unit for the xsphere slow control service.
#
# Install:
#   sudo cp xsphere-slowcontrol.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable xsphere-slowcontrol
#   sudo systemctl start xsphere-slowcontrol
#
# Logs:
#   journalctl -u xsphere-slowcontrol -f

[Unit]
Description=xsphere Slow Control Service
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=xbox
# Must run from the repo root so the `slowcontrol` package is importable.
WorkingDirectory=/home/xbox/xsphere-slow-control
ExecStart=/home/xbox/miniconda3/envs/slowcontrol/bin/python -u -m slowcontrol.app -c slowcontrol/config.yaml
Restart=on-failure
RestartSec=60
# Cap the restart loop: at most 5 starts per 10 minutes, then go to "failed".
# This pairs with the MQTT self-watchdog in the service heartbeat loop
# (slowcontrol/core/service.py): if the watchdog trips, the service exits
# non-zero, and systemd restarts it 60 s later. If something is genuinely
# broken (broker hard down, hardware failure) the service goes to "failed"
# after ~10 minutes of futile retries and a human is paged.
StartLimitIntervalSec=600
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
